# angel-api.service
[Unit]
Description=Angel One Order Handler API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/usr/bin/python3 /home/ubuntu/angel_order_handler.py

# Angel One Credentials - SET YOUR CREDENTIALS HERE
Environment="ANGEL_API_KEY=YOUR_API_KEY"
Environment="ANGEL_CLIENT_ID=YOUR_CLIENT_ID"
Environment="ANGEL_PASSWORD=YOUR_PASSWORD"
Environment="ANGEL_TOTP_SECRET=YOUR_TOTP_SECRET"

# Other environment
Environment="PYTHONUNBUFFERED=1"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=angel-api

# Auto-restart on failure
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
