# live-price-updater.service
[Unit]
Description=Stock Yard Live Price Updater
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Stock-Yard-Public
ExecStart=/usr/bin/python3 /home/ubuntu/Stock-Yard-Public/live_price_updater.py --data-file ./data.json --interval 10
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables (if needed)
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target
