# telegram-webhook.service
[Unit]
Description=Stock Yard Telegram Webhook
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
Environment="TELEGRAM_BOT_TOKEN=8253327701:AAGNFzBJ8QwKw8x8Hg-tlvWHg18DD4lgogQ"
Environment="TELEGRAM_CHAT_ID=8901309420"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/ubuntu"
ExecStart=/usr/bin/python3 /home/ubuntu/telegram_webhook_handler.py
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-webhook

[Install]
WantedBy=multi-user.target
